# kasa-monitor.conf
[program:kasa-monitor]
command=poetry run python main.py
directory=/home/c8763yee/kasa-monitor
autostart=true
user=c8763yee
