# mqtt2db/mqtt2db.conf
[program:mqtt2db]
directory=/home/c8763yee/mqtt2db
autostart=true
autorestart=true
startretries=10
user=c8763yee
logfile_maxbytes=0
logfile_backups=0
command=poetry run python3 app.py
